# deploy/systemd/cansat-radio-protocol.service
# CanSat — radio draad-protocol op de Zero 2 W (autostart, geen SSH nodig).
#
# Installatie (één keer op de Zero, paden aanpassen naar jullie user + clone):
#   sudo cp deploy/systemd/cansat-radio-protocol.service /etc/systemd/system/
#   sudo nano /etc/systemd/system/cansat-radio-protocol.service   # WorkingDirectory + ExecStart
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now cansat-radio-protocol.service
#
# Vereisten:
#   - SPI aan; venv met pip install -e . (en evt. .[rpi])
#   - Timezone (bv. Brussel): sudo timedatectl set-timezone Europe/Brussels
#   - SET TIME over de radio: dit service draait als **root** zodat de klok gezet kan worden
#     (clock_settime). Zonder root: timedatectl + sudoers, of weiger SET TIME (ERR TIME).
#
# Stoppen: STOP RADIO vanaf de Pico, of: sudo systemctl stop cansat-radio-protocol.service

[Unit]
Description=CanSat RFM69 wire protocol (CONFIG/MISSION listener)
After=network-pre.target local-fs.target
Wants=network-pre.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/icw/cansat_mission_2026
Environment=PYTHONUNBUFFERED=1
Environment=GPIOZERO_PIN_FACTORY=rpigpio
# Zorg dat de photo/log-map bestaat vóór de service start; anders faalt de
# preflight-check ``LOG`` (de service draait als root, dus ``~/photos`` ≠
# ``/home/icw/photos``) — we maken daarom de juiste user-map en chmod'en die.
ExecStartPre=/bin/mkdir -p /home/icw/photos
ExecStartPre=/bin/chown icw:icw /home/icw/photos
# Idem voor de binary log-map (cansat_continuous.bin + per-sessie *.bin).
# Zonder dit zou Path.home() onder root resolveren naar /root/cansat_logs/
# en zou icw die files niet zonder sudo kunnen lezen.
ExecStartPre=/bin/mkdir -p /home/icw/cansat_logs
ExecStartPre=/bin/chown icw:icw /home/icw/cansat_logs
# Pas --freq / --node / --dest / pins aan jullie hardware aan.
# --freq is enkel fallback: zodra de Pico ``SET FREQ`` heeft gedaan, staat
# ``config/radio_runtime.json`` op beide kanten in sync.
ExecStart=/home/icw/cansat_mission_2026/.venv/bin/python scripts/cansat_radio_protocol.py --verbose --poll 1.0 --reply-delay 0.08 --photo-dir /home/icw/photos --log-dir /home/icw/cansat_logs --bme280-addr 0x77
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
